# deploy/systemd-release/blueprint-re-frontend.service
[Unit]
Description=Blueprint RE Frontend
After=network.target blueprint-re-backend.service
Requires=blueprint-re-backend.service

[Service]
Type=simple
EnvironmentFile=%h/.config/blueprint-re/frontend.env
WorkingDirectory=__RELEASE_ROOT__/frontend-standalone/frontend
ExecStart=/bin/sh -c 'HOSTNAME=${FRONTEND_HOST} PORT=${FRONTEND_PORT} exec __NODE_BIN__ server.js'
Restart=on-failure
RestartSec=2
TimeoutStopSec=10
KillMode=control-group
SendSIGKILL=yes
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=default.target
