# deploy/systemd-release/blueprint-re-backend.service
[Unit]
Description=Blueprint RE Backend
After=network.target blueprint-re-manager-agent.service

[Service]
Type=simple
EnvironmentFile=%h/.config/blueprint-re/backend.env
WorkingDirectory=__RELEASE_ROOT__
ExecStart=__ENV__/bin/uvicorn app.main:app --host ${BACKEND_HOST} --port ${BACKEND_PORT}
Restart=on-failure
RestartSec=2
TimeoutStopSec=10
KillMode=control-group
SendSIGKILL=yes
# NoNewPrivileges is intentionally omitted for the backend service because it
# spawns executor sandbox processes via bubblewrap, which may need setuid or
# unprivileged userns privileges that NoNewPrivileges would block. The bwrap
# compatibility check is run at install time; if the host requires a setuid
# bwrap, backend service hardening must not prevent it from functioning.
PrivateTmp=yes

[Install]
WantedBy=default.target
